# src/sysprofd/sysprof3.service.in
[Unit]
Description=Sysprof Daemon

[Service]
Type=dbus
BusName=org.gnome.Sysprof3
ExecStart=@sysprofdprivdir@/sysprofd

